# etc/s6/lala/finish
#!/bin/sh

echo "Oops, lala died. Check logs at /var/log/lala"
