# etc/s6/lala/run
#!/bin/sh

mkdir -p /var/log
echo "Spawning lala"

exec /etc/bin/my-binary \
  >>/var/log/lala 2>&1

