# etc/s6/lulu/finish
#!/bin/sh

echo "Oops, lulu died. Check logs at /var/log/lulu"
