# etc/s6/.s6-svscan/finish
#!/bin/sh

echo "finishing whole thing"
