# etc/s6/.s6-svscan/finish
#!/bin/sh

# This script is run once a termination happens on pid1.
# Just acknowledge that that is happening.
echo "finishing whole thing"
